# valheim/valheim.crontab
0 6 * * * systemctl restart valheim
*/10 * * * * /opt/patched-dockerfiles/valheim/backup.sh /var/valheim/backup-10min
0 5 * * * /opt/patched-dockerfiles/valheim/backup.sh /var/valheim/backup-daily
10 5 * * 0 /opt/patched-dockerfiles/valheim/backup.sh /var/valheim/backup-weekly
