# minecraft/ts3server.crontab
15 6 * * * systemctl restart minecraftserver